# jozi.conf
[program:jozi]
command = /srv/jozi/bin/gunicorn_start.bash
user = root
stdout_logfile = /srv/jozi/log/gunicorn_supervisor.log
redirect_stderr = true
